# files/init.d/cumtlogin
#!/bin/sh /etc/rc.common

# 启动顺序
START=95

start() {
    config_load cumtlogin

    local enable=$(uci get cumtlogin.userconfig.enable)
    if [ "$enable" = "0" ]; then
        echo "cumtlogin is disabled."
        exit 0
    fi

    local testtime=$(uci get cumtlogin.userconfig.testtime)

    # Run login immediately on start
    /usr/sbin/cumtlogin >>/tmp/cumtloginlog 2>&1

    # Update crontab entry
    if grep -q "/usr/sbin/cumtlogin" /etc/crontabs/root; then
        sed -i '/\/usr\/sbin\/cumtlogin/d' /etc/crontabs/root
    fi
    echo "* */$testtime * * * /usr/sbin/cumtlogin >>/tmp/cumtloginlog 2>&1" >> /etc/crontabs/root
    /etc/init.d/cron reload
}

stop() {
    # Remove crontab entry and reload cron
    if grep -q "/usr/sbin/cumtlogin" /etc/crontabs/root; then
        sed -i '/\/usr\/sbin\/cumtlogin/d' /etc/crontabs/root
        /etc/init.d/cron reload
    fi
}

restart() {
    stop
    start
}
